# root/etc/periodic/15min/cleanup_sockets
#!/bin/sh

# Remove old sockets
find /tunnel -type s | xargs -I {} sh -c 'socat -u OPEN:/dev/null UNIX:{} || rm {}'
